# docker/php/program-apache.conf
[program:apache]
command = docker-php-entrypoint apache2-foreground
